# systemd/claude-feishu.service.tmpl
[Unit]
Description=Claude Code Feishu channel daemon
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/env bash -c 'cd ${PLUGIN_ROOT} && bun src/daemon.ts'
Restart=on-failure
RestartSec=5
Environment=FEISHU_STATE_DIR=${HOME}/.claude/channels/feishu
# Ensure bun and claude are on PATH (systemd user env is minimal):
Environment=PATH=${HOME}/.local/bin:${HOME}/.bun/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
