# launchd/com.claude-feishu.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-feishu</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd ${PLUGIN_ROOT} &amp;&amp; exec bun src/daemon.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>${PLUGIN_ROOT}</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <!-- launchd analog of the systemd unit's KillMode=process. Without this,
       `launchctl bootout` SIGTERMs the daemon's entire process group, which
       includes the tmux server spawned via ensureTmuxSession and every
       claude child underneath it — killing active sessions on every
       restart and defeating shim reconnect. -->
  <key>AbandonProcessGroup</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>${HOME}</string>
    <key>FEISHU_STATE_DIR</key>
    <string>${HOME}/.claude/channels/feishu</string>
    <!-- launchd user-agent env is minimal; enumerate PATH so bun/tmux/claude resolve. -->
    <key>PATH</key>
    <string>${HOME}/.local/bin:${HOME}/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>${HOME}/.claude/channels/feishu/daemon.log</string>
  <key>StandardErrorPath</key>
  <string>${HOME}/.claude/channels/feishu/daemon.log</string>
</dict>
</plist>
